# deploy/vellum.service
[Unit]
Description=Vellum backend (FastAPI + web, localhost only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vellum
Group=vellum
WorkingDirectory=/opt/vellum/api
Environment=VELLUM_WEB_DIST=/opt/vellum/web/dist
# Change VELLUM_PORT if 18080 is taken (e.g. another app on the box); keep the
# tunnel's remote port in sync. The app stays bound to 127.0.0.1 either way.
Environment=VELLUM_PORT=18080
ExecStart=/opt/vellum/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port ${VELLUM_PORT} --env-file .env
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
