# deploy/vellum.service
[Unit]
Description=Vellum backend (FastAPI + web, localhost only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vellum
Group=vellum
WorkingDirectory=/opt/vellum/api
Environment=VELLUM_WEB_DIST=/opt/vellum/web/dist
ExecStart=/opt/vellum/api/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 18080 --env-file .env
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
